# systemd/frigate-sync.service
[Unit]
Description=Frigate Recording Sync to Google Drive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=daniel
WorkingDirectory=/home/daniel/frigate-storage
ExecStart=/home/daniel/frigate-storage/sync-recordings.sh
StandardOutput=append:/home/daniel/frigate-storage/sync.log
StandardError=append:/home/daniel/frigate-storage/sync.log

# Timeout after 4 hours (large uploads may take time)
TimeoutStartSec=14400

[Install]
WantedBy=multi-user.target
